# school.service
[Unit]
Description=School Flask Application (Student Management)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/school
Environment="PATH=/var/www/school/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FLASK_ENV=production"
Environment="PORT=11000"
Environment="SECRET_KEY=3840567e0b897177054785797b895fb84cef6b64ac89ae68688045a10cdbb0e1"
Environment="DATABASE_URL=postgresql+psycopg2://school_user:your_db_password@localhost:5432/school_db"
ExecStart=/var/www/school/venv/bin/gunicorn --config gunicorn.conf.py wsgi:app

# Restart policy
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
